# sv/agent/run
#!/bin/bash

sv start zabbix || exit 1

source /etc/envvars
env

exec 2>&1
chpst -u zabbix zabbix_agentd -f
